# randomizer.service
[Unit]
Description=Server Randomizer

[Service]
ExecStart=python3 /usr/bin/mirrorlist-randomizer.py

[Install]
WantedBy=multi-user.target 
